# config/celeryd/shutdown_vm/shutdown_watcher.conf
start on startup

exec start wait_for_file FILE_PATH=/tmp/shutdown.dat
