# config/celeryd/shutdown_vm/shutdown_vm.conf
# /etc/init/myapp.conf
start on file FILE_PATH=/tmp/shutdown.dat

script
  rm /var/run/shutdown.dat
  shutdown -h +5 
end script
